# scripts/survey/systemd/backup.path
[Unit]
Description=Checks if log file path has been changed

[Path]
PathChanged=%h/log/partI/statistic

[Install]
WantedBy=default.target
